# config/supervisor.conf
; Notes:
; priority=1 --> Lower priorities indicate programs that start first and shut down last
; killasgroup=true --> send kill signal to child processes too

[program:frappe-framework-frappe-web]
command=/home/nandhakumar/frappe-framework/env/bin/gunicorn -b 127.0.0.1:8000 -w 13 --max-requests 5000 --max-requests-jitter 500 -t 120 frappe.app:application --preload
priority=4
autostart=true
autorestart=true
stdout_logfile=/home/nandhakumar/frappe-framework/logs/web.log
stderr_logfile=/home/nandhakumar/frappe-framework/logs/web.error.log
user=nandhakumar
directory=/home/nandhakumar/frappe-framework/sites


[program:frappe-framework-frappe-schedule]
command=/usr/local/bin/bench schedule
priority=3
autostart=true
autorestart=true
stdout_logfile=/home/nandhakumar/frappe-framework/logs/schedule.log
stderr_logfile=/home/nandhakumar/frappe-framework/logs/schedule.error.log
user=nandhakumar
directory=/home/nandhakumar/frappe-framework

[program:frappe-framework-frappe-default-worker]
command=/usr/local/bin/bench worker --queue default
priority=4
autostart=true
autorestart=true
stdout_logfile=/home/nandhakumar/frappe-framework/logs/worker.log
stderr_logfile=/home/nandhakumar/frappe-framework/logs/worker.error.log
user=nandhakumar
stopwaitsecs=1560
directory=/home/nandhakumar/frappe-framework
killasgroup=true
numprocs=1
process_name=%(program_name)s-%(process_num)d

[program:frappe-framework-frappe-short-worker]
command=/usr/local/bin/bench worker --queue short
priority=4
autostart=true
autorestart=true
stdout_logfile=/home/nandhakumar/frappe-framework/logs/worker.log
stderr_logfile=/home/nandhakumar/frappe-framework/logs/worker.error.log
user=nandhakumar
stopwaitsecs=360
directory=/home/nandhakumar/frappe-framework
killasgroup=true
numprocs=1
process_name=%(program_name)s-%(process_num)d

[program:frappe-framework-frappe-long-worker]
command=/usr/local/bin/bench worker --queue long
priority=4
autostart=true
autorestart=true
stdout_logfile=/home/nandhakumar/frappe-framework/logs/worker.log
stderr_logfile=/home/nandhakumar/frappe-framework/logs/worker.error.log
user=nandhakumar
stopwaitsecs=1560
directory=/home/nandhakumar/frappe-framework
killasgroup=true
numprocs=1
process_name=%(program_name)s-%(process_num)d






[program:frappe-framework-redis-cache]
command=/usr/bin/redis-server /home/nandhakumar/frappe-framework/config/redis_cache.conf
priority=1
autostart=true
autorestart=true
stdout_logfile=/home/nandhakumar/frappe-framework/logs/redis-cache.log
stderr_logfile=/home/nandhakumar/frappe-framework/logs/redis-cache.error.log
user=nandhakumar
directory=/home/nandhakumar/frappe-framework/sites

[program:frappe-framework-redis-queue]
command=/usr/bin/redis-server /home/nandhakumar/frappe-framework/config/redis_queue.conf
priority=1
autostart=true
autorestart=true
stdout_logfile=/home/nandhakumar/frappe-framework/logs/redis-queue.log
stderr_logfile=/home/nandhakumar/frappe-framework/logs/redis-queue.error.log
user=nandhakumar
directory=/home/nandhakumar/frappe-framework/sites



[program:frappe-framework-redis-socketio]
command=/usr/bin/redis-server /home/nandhakumar/frappe-framework/config/redis_socketio.conf
priority=1
autostart=true
autorestart=true
stdout_logfile=/home/nandhakumar/frappe-framework/logs/redis-socketio.log
stderr_logfile=/home/nandhakumar/frappe-framework/logs/redis-socketio.error.log
user=nandhakumar
directory=/home/nandhakumar/frappe-framework/sites



[program:frappe-framework-node-socketio]
command=/usr/bin/node /home/nandhakumar/frappe-framework/apps/frappe/socketio.js
priority=4
autostart=true
autorestart=true
stdout_logfile=/home/nandhakumar/frappe-framework/logs/node-socketio.log
stderr_logfile=/home/nandhakumar/frappe-framework/logs/node-socketio.error.log
user=nandhakumar
directory=/home/nandhakumar/frappe-framework


[group:frappe-framework-web]
programs=frappe-framework-frappe-web,frappe-framework-node-socketio



[group:frappe-framework-workers]
programs=frappe-framework-frappe-schedule,frappe-framework-frappe-default-worker,frappe-framework-frappe-short-worker,frappe-framework-frappe-long-worker




[group:frappe-framework-redis]
programs=frappe-framework-redis-cache,frappe-framework-redis-queue,frappe-framework-redis-socketio
